# deployment/com.service.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- launchd plist for the service. Update paths before installation. -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.porthandler</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/path/to/PortHandler/.venv/bin/python</string>
        <string>/path/to/PortHandler/src/main.py</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/path/to/PortHandler</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/tmp/porthandler.log</string>
    
    <key>StandardErrorPath</key>
    <string>/tmp/porthandler.error.log</string>
</dict>
</plist>